# etc/init.d/jerbild
#!/sbin/runscript
#
# INIT Script for the Jerbil - Ruby Object Broker


depend() {
  use net
  use logger
}

export RUBYOPT="rubygems"
export RUBYLIB=""

# create the options for starting the daemon from the settings in the
# /etc/conf.d/jerbild file

myopts=""
[ -n "${NO_DAEMON}" ] && myopts="${myopts} -n"
[ -n "${CONF_FILE}" ] && myopts="${myopts} -c ${CONF_FILE}"
[ "${VERBOSE}" == "true" ] && [ "${QUIET}" != "true" ] && myopts="${myopts} -V"
[ -n "${NO_SYSLOG}" ] && myopts="${myopts} -S"


start() {
  ebegin "Starting Jerbil Server"

    #/bin/su -m -c "/usr/sbin/jerbild ${myopts}" - jermine
    /usr/sbin/jerbild ${myopts}

  eend $?
}

status() {
  ebegin "Checking Jerbil Status"
  /usr/sbin/jerbil-status ${myopts}
  eend $?
}


stop() {
  ebegin "Stopping Jerbil Server"
    /usr/sbin/jerbil-stop ${myopts}
  eend $?
}

