# etc/init.d/jerbild
#!/sbin/runscript
#
# INIT Script for the Jerbil - Ruby Object Broker


depend() {
  use net
  use logger
}
  
describe() {
  einfo "Reliable Broker for Ruby Services"
}

export RUBYOPT="rubygems"
export RUBYLIB=""


# create the options for starting the daemon from the settings in the
# /etc/conf.d/jerbild file

myopts=""
[ -n "${NO_DAEMON}" ] && myopts="${myopts} -n"
[ -n "${CONF_FILE}" ] && myopts="${myopts} -c ${CONF_FILE}"
[ "${VERBOSE}" == "true" ] && [ "${QUIET}" != "true" ] && myopts="${myopts} -V"
[ -n "${NO_SYSLOG}" ] && myopts="${myopts} -S"


start() {
  ebegin "Starting Jerbil Server"
  
  # /var/run is now temporary and may need to be created at each boot
  #if [ ! -d /var/run/jerbil ] ; then
  #  mkdir /var/run/jerbil
  #  chown jerbil:jerbil /var/run/jerbil
  #  chmod 775 /var/run/jerbil
  #  einfo "Created /var/run/jerbil"
  #fi
  
  checkpath -q -d -m 0755 -o jerbil:jerbil /var/run/jerbil

  /bin/su -c "${SBIN_PATH}jerbild ${myopts}" - jerbil

  eend $?
}

status() {
  ebegin "Checking Jerbil Status"
    /bin/su -c "${SBIN_PATH}jerbil-status ${myopts}" - jerbil
  eend $?
}


stop() {
  ebegin "Stopping Jerbil Server"
    /bin/su -c "${SBIN_PATH}jerbil-stop ${myopts}" - jerbil
  eend $?
}

