# etc/init.d/jserviced
#!/sbin/runscript
#
# INIT Script for a Jerbil Service

# create dependencies from the list of services in $NEEDS and $USES
#depend() {
#  for need_service in $NEEDS
#  do
#    need $need_service
#  done
#  for use_service in $USES
#  do
#    use $use_service
#  done
#}

# make sure there are no unexpected rubylibs
export RUBYLIB=""
export RUBYOPT="rubygems"


myopts=""
[ -n "${NO_DAEMON}" ] && myopts="${myopts} -n"
[ -n "${CONF_FILE}" ] && myopts="${myopts} -c ${CONF_FILE}"
[ "${VERBOSE}" == "true" ] && [ "${QUIET}" != "true" ] && myopts="${myopts} -V"
[ -n "${NO_SYSLOG}" ] && myopts="${myopts} -S"
[ -n "${LOG_DAEMON}" ] && myopts="${myopts} -l"
[ "$SERVICE_NAME" == "" ] && SERVICE_NAME="${SVCNAME}"
[ "$SERVICE_USER" == "" ] && SERVICE_USER="jerbil"

start() {
  ebegin "Starting Jerbil Service ${SERVICE_NAME}"


    # create the options for starting the daemon from the settings in the
    # /etc/conf.d/jservice file (same name as this link)

    /bin/su -c "${SBIN_PATH}jserviced ${myopts} -s ${SERVICE_NAME}" - ${SERVICE_USER}

  eend $?
}

status() {
  ebegin "Checking status for Jerbil Service ${SERVICE_NAME}"
    /bin/su -c "${SBIN_PATH}jservice-status ${myopts} -s ${SERVICE_NAME}" -${SERVICE_USER}
  eend $?
}


stop() {
  ebegin "Stopping Jerbil Service ${SERVICE_NAME}"
    /bin/su -c "${SBIN_PATH}jservice-stop ${myopts} -s ${SERVICE_NAME}" - ${SERVICE_USER}
  eend $?
}
  
describe() {
  einfo ${DESCRIPTION}

}

