# init.d/jerbild
#!/sbin/runscript
#
# INIT Script for the Numbat System Controller

opts="${opts} install"
export RUBYOPT="-rauto_gem"
export RUBYLIB="/usr/local/lib"

start() {
  ebegin "Starting Jerbil Server"

    if [ ! -d /var/run/jermine ]; then
      install
    fi

    # create the options for starting the daemon from the settings in the
    # /etc/conf.d/jerbild file
    myopts=""
    [ -n "${NO_DAEMON}" ] && myopts="${myopts} -n"
    [ -n "${CONF_FILE}" ] && myopts="${myopts} -c ${CONF_FILE}"
    [ "${VERBOSE}" == "true" ] && [ "${QUIET}" != "true" ] && myopts="${myopts} -V"
    [ -n "${NO_SYSLOG}" ] && myopts="${myopts} -S"

    /bin/su -m -c "/usr/local/sbin/jerbild ${myopts}" - jermine

  eend $?
}

status() {
  ebegin "Checking Jerbil Status"
    /usr/local/bin/jerbil -q
  eend $?
}


stop() {
  ebegin "Stopping Jerbil Server"
    /usr/local/sbin/jerbil-stop
  eend $?
}

install() {
  einfo "Installing Jerbil Server"
  ebegin "Creating User"
    /usr/sbin/useradd -d /var/run/jermine -u 603 -s /bin/bash jermine >&/dev/null
  eend $?
  ebegin "Creating Group"
    /usr/sbin/groupadd -g 603 -f jermine >&/dev/null
    /usr/sbin/gpasswd -a jermine jermine >&/dev/null
    /usr/sbin/gpasswd -a robert jermine >&/dev/null
  ebegin "Creating directories"
    mkdir -p /var/log/jermine >&/dev/null
    mkdir -p /var/run/jermine >&/dev/null
  eend $?
  ebegin "Changing log directory ownership"
    chown jermine:jermine /var/log/jermine >&/dev/null
    chown jermine:jermine /var/run/jermine >&/dev/null
    chmod g+w /var/log/jermine >&/dev/null
    chmod g+w /var/run/jermine >&/dev/null
  eend $?
}
