# init.d/jerbild
#!/sbin/runscript
#
# INIT Script for the Numbat System Controller

opts="${opts} install"
export RUBYOPT="-rauto_gem"
export RUBYLIB="/usr/local/lib"

start() {
  ebegin "Starting Jerbil Server"

    if [ ! -d /var/run/jermine ]; then
      install
    fi

    # create the options for starting the daemon from the settings in the
    # /etc/conf.d/jerbild file
    myopts=""
    [ -n "${NO_DAEMON}" ] && myopts="${myopts} -n"
    [ -n "${CONF_FILE}" ] && myopts="${myopts} -c ${CONF_FILE}"
    [ "${VERBOSE}" == "true" ] && [ "${QUIET}" != "true" ] && myopts="${myopts} -V"
    [ -n "${NO_SYSLOG}" ] && myopts="${myopts} -S"

    #/bin/su -m -c "/usr/sbin/jerbild ${myopts}" - jermine
    /usr/sbin/jerbild ${myopts}

  eend $?
}

status() {
  ebegin "Checking Jerbil Status"
    /usr/local/bin/jerbil -q
  eend $?
}


stop() {
  ebegin "Stopping Jerbil Server"
    /usr/sbin/jerbil-stop
  eend $?
}

